# systemd/cert.timer
[Unit]
Description=Cert renewal

[Timer]
OnBootSec=5min
OnCalendar=weekly

[Install]
WantedBy=timers.target